# backend-go/config/bot.ini
[program:bot]
command=/Users/sh001/Documents/codes/ivan/liquidation/backend-go/bin/liquidation-bot
autostart=true
autorestart=true
stderr_logfile=/Users/sh001/Documents/codes/ivan/liquidation/backend-go/logs/bot.err.log
stdout_logfile=/Users/sh001/Documents/codes/ivan/liquidation/backend-go/logs/bot.out.log
